# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

crontab -u husky /etc/crontabs/husky
cron -f -l